# apps/desktop/services/optimus-vault-events/optimus-vault-events.service
[Unit]
Description=Optimus BotVault authoritative note event feed (:9128)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/optimus-vault-events/server.py
Restart=always
RestartSec=3
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
